# infrastructure/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:caddy]
command=/usr/sbin/caddy

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F